# packaging/systemd/spong-sgt-sync.timer
[Unit]
Description=Spong竊担GT link reconciliation (cada 5 min)

[Timer]
# Primer disparo 2 min despuﾃｩs de boot, luego cada 5 min.
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
