# packaging/spong-client/DEBIAN/prerm
#!/bin/bash
set -e

echo "Deteniendo spong-client..."
systemctl stop    spong-client 2>/dev/null || true
systemctl disable spong-client 2>/dev/null || true
